# sub/avm-eac_rails_base1/template/avm/eac_rails_base1/systemd_unit/tasks_scheduler.service
[Unit]
Description='%%description%%'
After=postgresql.service

[Service]
Type=forking
%%service_exec_lines%%
Restart=%%restart%%
User=%%user%%

[Install]
WantedBy=multi-user.target
